# src/Systemd/nanohatoled.service.in
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=Nanopi Oled Hat Software
After=network.target remote-fs.target syslog.target time.target


[Service]
PIDFile=/run/nanohatoled.pid
ExecStartPre=-/bin/touch /var/log/nanohatoled.log
ExecStart=/bin/sh -c '@BIN_INSTALL_DIR@/nanohatoled'
ExecReload=/bin/kill -s HUP $MAINPID
WorkingDirectory=@BIN_INSTALL_DIR@
StandardOutput=inherit
StandardError=inherit
Restart=always

[Install]
WantedBy=multi-user.target
